# jjp-agent/plists/com.jjp.briefing.sunday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jjp.briefing.sunday</string>

    <key>ProgramArguments</key>
    <array>
        <!-- UPDATED BY SETUP SCRIPT -->
        <string>NODE_PATH</string>
        <string>src/briefing-standalone.js</string>
        <string>weekly</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- UPDATED BY SETUP SCRIPT -->
    <string>AGENT_PATH</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/jjp-briefing-sunday.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/jjp-briefing-sunday.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
